# launchd/com.instantvim.qt.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.instantvim.qt</string>

  <!--
    Launching the raw Ghostty binary directly (rather than via `open`) is
    intentional: launchd's KeepAlive needs to track the actual long-lived
    process, and `open -na` would just exit immediately after handing off,
    leaving launchd nothing to supervise. Empirically verified working on
    Ghostty 1.3.1/macOS: this starts a real, fully functional hidden
    instance whose quick terminal responds to its own keybind.

    __QT_CONFIG_PATH__ is substituted by install.sh with the absolute path
    to the generated quick-terminal.config (see host/quick-terminal.config).
  -->
  <key>ProgramArguments</key>
  <array>
    <string>/Applications/Ghostty.app/Contents/MacOS/ghostty</string>
    <string>--config-default-files=false</string>
    <string>--config-file=__QT_CONFIG_PATH__</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProcessType</key>
  <string>Interactive</string>
  <key>LimitLoadToSessionType</key>
  <string>Aqua</string>

  <key>StandardOutPath</key>
  <string>/tmp/instantvim-qt.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/instantvim-qt.log</string>
</dict>
</plist>
